# wg-registration@.service
[Unit]
Description=Wireguard peer registration system

[Service]
User=%i
ExecStart=/usr/bin/env wireguard_registration_listener.sh

[Install]
WantedBy=multi-user.target
